# packaging/org.tizen.account.manager.service
[Unit]
Description=Account Service D-BUS

[D-BUS Service]
Name=org.tizen.account.manager
Exec=/usr/bin/account-svcd
